# templates/default/supervisord.sentinal-master.conf.erb
[program:sentinal_master]
command = /usr/bin/python /var/sentinal-master.py
autostart=true
autorestart=true

stopsignal=KILL
stopasgroup = true
killasgroup = true

stdout_logfile=/tmp/sentinal_master.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1

stderr_logfile = /tmp/sentinal_master.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1

[group:sentinal_master_server]
programs=sentinal_master

