# templates/default/supervisord.sentinal.conf.erb
[program:chef-client]
command = chef-client -i 240 -s 0 -f -N <%=node.name%>
autostart=true
autorestart=true

stopsignal=KILL
stopasgroup = true
killasgroup = true

stdout_logfile=/tmp/chef-client.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

stderr_logfile = /tmp/chef-client.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[group:chef-client_server]
programs=chef-client
