# templates/default/supervisord.sentinal.conf.erb
[program:sentinal]
command = /var/redis-3.0.0/src/redis-sentinel /var/sentinal.conf
autostart=true
autorestart=true

stopsignal=KILL
stopasgroup = true
killasgroup = true

stdout_logfile=/tmp/sentinal.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1

stderr_logfile = /tmp/sentinal.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1

[group:sentinal_server]
programs=sentinal

